# weather_report/weather_report
### BEGIN INIT INFO
# Provides: weather_report
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start Fluentd + weather_report at boot time
# Description: Start Fluentd + weather_report at boot time.
### END INIT INFO
 
#! /bin/sh
# /etc/init.d/weather_report
 
USER=pi
HOME=/home/pi
 
export USER HOME
 
case "$1" in
    start)
        echo "Starting weather report"
	cd /home/pi/raspi/weather_report
	nohup /usr/local/bin/fluentd -c fluentd.conf &
	nohup sudo python weather_report.py & 
        ;;
 
    stop)
        echo "Please reboot RasPi to restart"
        ;;
 
    *)
        echo "Usage: /etc/init.d/wether_report {start|stop}"
        exit 1
        ;;
esac
 
exit 0
